# tools/openrc/atlasdaemon
#!/sbin/runscript
# Short-Description: ATLAS manager deamon
# Description:       ATLAS manager daemon used to start the server and keep it updated

source /etc/atlasmanager/atlasmanager.cfg

NAME="ShooterGameServer"
LOGFILE="${logdir}/${NAME}.log"
DAEMON="/usr/bin/atlasmanager"

depend(){
    need net
}

start(){
    INSTANCE="${RC_SVCNAME#*.}"
    ebegin "Starting ATLAS manager daemon"
    ulimit -n 100000
    if [ "$INSTANCE" != "$RC_SVCNAME" ]; then
        "$DAEMON" start "@${INSTANCE}"
    else
        "$DAEMON" start "@all"
    fi
    eend $?
}

stop(){
    INSTANCE="${RC_SVCNAME#*.}"
    ebegin "Stopping ATLAS manager daemon"
    if [ "$INSTANCE" != "$RC_SVCNAME" ]; then
        "$DAEMON" stop "@${INSTANCE}"
    else
        "$DAEMON" stop "@all"
    fi
    eend $?
}

